# installer/spawnwp-telemetry.service
[Unit]
Description=Send optional aggregate SpawnWP telemetry

[Service]
Type=oneshot
ExecStart=/usr/local/lib/spawnwp/telemetry.py send heartbeat
