# installer/spawnwp-site-expiry.service
[Unit]
Description=SpawnWP destruction of expired temporary sites

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /srv/wp-dev/scripts/site-expiry.sh
