# installer/spawnwp-image-gc.service
[Unit]
Description=SpawnWP optional auto-delete of unused PHP images

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /srv/wp-dev/scripts/image-gc.sh
