# installer/wp-cockpit.service
[Unit]
Description=SpawnWP cockpit
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/srv/wp-cockpit
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/srv/wp-cockpit/venv/bin/uvicorn app:app --host 127.0.0.1 --port 9393 --workers 1
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
